# packages/java_stuff/java/profile.d/java.profile
# needs to point where jre folder exists
export JAVA_HOME=/storage/.kodi/addons/program.java

# look for libbluray-*.jar in this folder (need / at the end)
export LIBBLURAY_CP=/storage/.kodi/addons/program.java/
# or set file directly
# export LIBBLURAY_CP=/storage/.kodi/addons/program.java/libbluray-j2se-1.0.2.jar
