# packages/sysutils/jre/profile.d/jre.profile
# needs to point where jre folder is set
export JAVA_HOME=/storage/.kodi/addons/program.jre/

# look for libbluray-*.jar in this folder (with / at the end)
export LIBBLURAY_CP=/storage/.kodi/addons/program.jre/
# or set file directly
# export LIBBLURAY_CP=/storage/.kodi/addons/program.jre/libbluray-j2se-1.0.1.jar

